# deploy/openai-account-control.service
[Unit]
Description=OpenAI Account Control Center
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/openai-account-control
EnvironmentFile=/opt/openai-account-control/.env
ExecStart=/opt/openai-account-control/.venv/bin/python -m uvicorn backend.app.main:app --host ${APP_HOST} --port ${APP_PORT}
Restart=on-failure
RestartSec=5
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
